# custom/cumbox2/systemd/cumbox2-fan.service
[Unit]
Description=CumeBox2 Fan Auto Control
After=multi-user.target sys-module-gpio-tee.service
Requires=sys-module-gpio-tee.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/cumbox2/scripts/fan_auto_control.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cumbox2-fan

# 环境变量
Environment="LC_ALL=C"

[Install]
WantedBy=multi-user.target