# custom/cumbox2/systemd/zram.service
[Unit]
Description=ZRAM Memory Compression
After=multi-user.target
Before=swap.target

[Service]
Type=oneshot
ExecStart=/usr/local/cumbox2/scripts/setup_zram.sh
ExecStop=/sbin/swapoff /dev/zram0 2>/dev/null || true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
